# deploy/blog-collector.service
[Unit]
Description=BLog 전주 버스 실시간 데이터 수집기
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=jiho
Group=blog
WorkingDirectory=/home/jiho/BLog
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/jiho/miniconda3/envs/Blog/bin/python -m src.collector
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
